# ubuntu-xboxdrv-20140415/debian/ubuntu-xboxdrv/etc/init/xboxdrv.conf
# xboxdrv.conf
start on filesystem
pre-start script
    # clean up joysticks
    rm -f /dev/input/js*
end script
exec xboxdrv --daemon --detach-kernel-driver --silent --dbus disabled
